# apps/mastodon/crontab
*/1 * * * * /app/bin/sidekiq-healthcheck
